# image-toolkit/misc/rc.local
#!/bin/bash

source /opt/vnoi/config.sh

# vim: ft=sh ts=4 noet
