# src/misc/rc.local
#!/bin/bash

source /opt/vnoi/config.sh

/opt/vnoi/bin/gnome-autostart.sh

# vim: ft=sh ts=4 noet
